# aws/packer/profile.d/conda.sh
#
# conda.sh:
#   Setup conda environment variables
#

PATH=$PATH:/usr/local/env/conda/bin

export PATH