# aws/packer/profile.d/custom-aliases.sh
#
# custom-aliases.sh:
#   Add custom aliases to all users
#

# slurm squeue, show wide display
alias sqq='squeue -o "%.18i %.9P %.50j %.8u %.2t %.10M %.6D %R"'